# evaluation/ansible/templates/wasimoff.service.j2
# {{ ansible_managed }}

[Unit]
Description=Project Wasimoff Broker
StartLimitIntervalSec=0

[Service]
Type=exec
{% for env in wasimoff_environment %}
Environment={{ env }}
{% endfor %}
ExecStart={{ wasimoff_binary }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target